# build/usr/lib/systemd/system/casaos-gateway.service
[Unit]
After=network.target
ConditionFileNotEmpty=/etc/casaos/gateway.ini
Description=CasaOS Gateway

[Service]
ExecStart=/usr/bin/casaos-gateway
PIDFile=/var/run/casaos/gateway.pid
#Restart=always

[Install]
WantedBy=multi-user.target
