# build/sysroot/usr/lib/systemd/system/casaos-gateway.service
[Unit]
After=network.target
ConditionFileNotEmpty=/etc/casaos/gateway.ini
Description=CasaOS Gateway

[Service]
ExecStartPre=/usr/bin/casaos-gateway -v
ExecStart=/usr/bin/casaos-gateway
PIDFile=/var/run/casaos/gateway.pid
Restart=always
Type=notify

[Install]
WantedBy=multi-user.target
