# build/usr/lib/casaos/casaos-gateway.service
[Unit]
Description=CasaOS Gateway
After=network.target
ConditionFileNotEmpty=/etc/casaos/gateway.ini

[Service]
EnvironmentFile=/etc/default/casaos-gateway
ExecStart=/usr/local/bin/

# See https://blog.devgenius.io/write-systemd-service-file-in-linux-centos-and-ubuntu-6accf1dc1109